# tools/launchd/com.fadacai.briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fadacai.briefing</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/YOUR_USERNAME/path/to/portfolio/tools/briefing_runner.sh</string>
  </array>

  <!-- Multiple ET trigger points. TZ env below interprets Hour in New York time
       (auto DST). A laptop often sleeps through any single fixed time, so we fire
       at several; briefing_runner.sh bails immediately if already sent today, and
       send_briefing.py dedups — so at most ONE briefing goes out per day. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key><string>America/New_York</string>
    <!-- Add PATH so launchd can find python3 and claude CLI -->
    <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/YOUR_USERNAME/.local/bin</string>
    <key>HOME</key><string>/Users/YOUR_USERNAME</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/path/to/portfolio/briefing-out/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/path/to/portfolio/briefing-out/launchd.err</string>

  <!-- Retry once if exit code is non-zero -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <!-- Wake the Mac from sleep at scheduled time (requires AC power) -->
  <key>WakeForNetwork</key>
  <true/>
</dict>
</plist>
